# roles/client/templates/offsitebuddy-check.timer.j2
[Unit]
Description=Run OffsiteBuddy restic check job {{ job.name }}

[Timer]
OnCalendar={{ job.check.schedule | default('weekly') }}
Persistent=true

[Install]
WantedBy=timers.target
